# deploy/com.trendfit.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.trendfit.dashboard</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/brunoliberatogirardi/Downloads/Dev/Trendfit/.venv/bin/python</string>
    <string>/Users/brunoliberatogirardi/Downloads/Dev/Trendfit/scripts/dashboard.py</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/brunoliberatogirardi/Downloads/Dev/Trendfit</string>
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>21</integer><key>Minute</key><integer>0</integer></dict>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/trendfit_dashboard.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/trendfit_dashboard.log</string>
</dict>
</plist>
